# deploy/agent-wiki-sweeper.service
[Unit]
Description=agent-wiki — daily sweep + reindex (oneshot)
After=network-online.target
Wants=network-online.target

# User-mode service: runs in your user systemd session via `systemctl --user`,
# so it can write freely under your project tree without the
# ProtectHome=read-only sandbox the main agent-wiki.service uses.
#
# Edit ExecStart to point at your checkout's deploy/run-daily.sh.

[Service]
Type=oneshot
ExecStart=%h/agent-wiki/deploy/run-daily.sh
# Anthropic API calls + sequential sweeps can take 5–10 minutes.
TimeoutStartSec=20min

[Install]
WantedBy=default.target
