# deploy/agent-wiki.service
[Unit]
# Edit User, Group, WorkingDirectory, EnvironmentFile, and ReadWritePaths to
# match your install. The paths below assume a checkout at /opt/agent-wiki
# owned by user `agent-wiki` — substitute your own.
Description=agent-wiki — read-only HTTP API + UI for project AGENTS.md docs
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=agent-wiki
Group=agent-wiki
WorkingDirectory=/opt/agent-wiki/service
EnvironmentFile=/opt/agent-wiki/service/.env
# tsx is fine for dev. For production, prefer the compiled output:
#   (cd service && npm run build) and then:
ExecStart=/usr/bin/node dist/server/index.js
Restart=on-failure
RestartSec=3
# Modest resource cap — service is a tiny Express app reading JSON.
MemoryMax=256M

# Hardening.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
# Allow writing the file-based session store under WorkingDirectory/.sessions.
ReadWritePaths=/opt/agent-wiki/service/.sessions

[Install]
WantedBy=multi-user.target
