# deploy/agent-wiki.service
[Unit]
Description=agent-wiki — read-only HTTP API + UI for project AGENTS.md docs
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=secorp
Group=secorp
WorkingDirectory=/home/secorp/termag/projects/agent-wiki/service
EnvironmentFile=/home/secorp/termag/projects/agent-wiki/service/.env
# tsx is fine for dev. For production, prefer the compiled output:
#   (cd service && npm run build) and then:
ExecStart=/usr/bin/node dist/server/index.js
Restart=on-failure
RestartSec=3
# Modest resource cap — service is a tiny Express app reading JSON.
MemoryMax=256M

# Hardening (mirrors what termag/rssreader use).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
# Allow writing the file-based session store under WorkingDirectory/.sessions.
ReadWritePaths=/home/secorp/termag/projects/agent-wiki/service/.sessions

[Install]
WantedBy=multi-user.target
